# template/parts/supervisor/supervisor.conf.in
[inet_http_server]
port = ${config:supervisord-port}
username = 
password = 

[supervisord]
logfile = ${buildout:directory}/parts/supervisor/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = ${buildout:directory}/var/supervisord.pid
nodaemon = false

[supervisorctl]
serverurl = http://localhost:${config:supervisord-port}

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface


[program:mysqld]
command = ${buildout:parts-directory}/mysql/libexec/mysqld
process_name = mysqld
directory = ${buildout:parts-directory}/mysql/libexec
priority = 10
redirect_stderr = false



[program:apache]
command = ${buildout:parts-directory}/apache/bin/httpd -DFOREGROUND
process_name = apache
directory = ${buildout:parts-directory}/apache/bin
priority = 10
redirect_stderr = false



[program:memcached]
command = ${buildout:parts-directory}/memcached/bin/memcached -v -p ${config:memcached-port} -U ${config:memcached-port}
process_name = memcached
directory = ${buildout:parts-directory}/memcached/bin
priority = 10
redirect_stderr = false




